# node/crontab
* *     * * *   root    /root/scripts/prevent_hanging_gpslogger.sh > /dev/null 2>&1
